# systemd/.config/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver

[Service]
Type=simple
ExecStart=/usr/bin/xscreensaver --no-splash

[Install]
WantedBy=graphical.target
