# systemd/.config/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver

[Service]
Type=simple
ExecStart=/usr/bin/xscreensaver -nosplash

[Install]
WantedBy=user-graphical-login.target
